# service/io.opencode.ntfy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.opencode.ntfy</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/usr/local/opt/opencode-ntfy/libexec/server.js</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/usr/local/var/log/opencode-ntfy.log</string>
    
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/opencode-ntfy.log</string>
    
    <key>WorkingDirectory</key>
    <string>/usr/local/opt/opencode-ntfy/libexec</string>
</dict>
</plist>
